# sparql/sparql.socket

[Unit]
Description=SPARQL socket

[Socket]
ListenStream=8081
Accept=yes

[Install]
WantedBy=sockets.target

